# src/build/modules/nodejs/files/etc/supervisor/conf.d/nodejs.conf
[program:nodejs]
command = /usr/bin/npm start /nodejs/data
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
